# scripts/m4bot-web.service
[Unit]
Description=M4Bot Web Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
User=m4bot
Group=m4bot
WorkingDirectory=/opt/m4bot
ExecStart=/opt/m4bot/venv/bin/python -m gunicorn "web.app:app" -k uvicorn.workers.UvicornWorker -b 0.0.0.0:5000 --log-level info
Restart=on-failure
Environment="PYTHONPATH=/opt/m4bot"
Environment="QUART_APP=web.app:app"
Environment="QUART_ENV=production"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target 